# scripts/cassettenest.conf
[program:cassettenest]
command=/usr/local/bin/cassettenest.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/cassettenest.err.log
stdout_logfile=/var/log/cassettenest.out.log
